# system/ovos.service
[Unit]
Description=OpenVoiceOS Software stack.

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=default.target
